# systemd/rct-dashboard.service
[Unit]
Description=RCT Power Dashboard (FastAPI)
After=network.target rct-poller.service

[Service]
Type=simple
User=solar
WorkingDirectory=/opt/rct-dashboard
ExecStart=/opt/rct-dashboard/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
